# system/backup/backup.cron
#locate in the /etc/cron.d/backup
30 19 * * 0 sh /backup/localBackup
